# init/qdiskd.service.in
[Unit]
Description=Qdiskd Service
ConditionKernelCommandLine=!nocluster
PartOf=corosync.service
After=corosync.service iscsi.service


[Service]
EnvironmentFile=-@INITCONFIGDIR@/corosync-qdiskd
ExecStart=@SBINDIR@/qdiskd $COROSYNC_QDISKD_OPTIONS
User=root
WorkingDirectory=/
StandardOutput=journal
StandardError=journal

Type=notify
Restart=on-abnormal
#WatchdogSec=9

[Install]
WantedBy=corosync.service
